# files/system/etc/profile.d/open.sh
# xdg-open alias — macOS-like "open" command
# Redirects stderr so GTK warnings don't pollute the terminal
alias open="xdg-open &>/dev/null"
