# files/system/etc/xdg/autostart/querencia-welcome.desktop
[Desktop Entry]
Type=Application
Name=Welcome to Querencia Linux
Name[de]=Willkommen bei Querencia Linux
Name[fr]=Bienvenue sur Querencia Linux
Name[es]=Bienvenido a Querencia Linux
Name[it]=Benvenuto in Querencia Linux
Name[pt]=Bem-vindo ao Querencia Linux
Name[nl]=Welkom bij Querencia Linux
Name[pl]=Witamy w Querencia Linux
Name[ru]=Добро пожаловать в Querencia Linux
Name[ja]=Querencia Linux へようこそ
Name[zh]=欢迎使用 Querencia Linux
Name[ko]=Querencia Linux에 오신 것을 환영합니다
Comment=Welcome screen for Querencia Linux
Comment[de]=Willkommensbildschirm für Querencia Linux
Comment[fr]=Écran de bienvenue pour Querencia Linux
Comment[es]=Pantalla de bienvenida de Querencia Linux
Comment[it]=Schermata di benvenuto di Querencia Linux
Comment[pt]=Ecrã de boas-vindas do Querencia Linux
Comment[nl]=Welkomstscherm voor Querencia Linux
Comment[pl]=Ekran powitalny Querencia Linux
Comment[ru]=Экран приветствия Querencia Linux
Comment[ja]=Querencia Linux ようこそ画面
Comment[zh]=Querencia Linux 欢迎界面
Comment[ko]=Querencia Linux 시작 화면
Exec=querencia-welcome-launcher
Icon=help-about
Terminal=false
NoDisplay=true
X-MATE-Autostart-enabled=true
X-MATE-Autostart-Delay=3
